# infra/deploy/bd-deploy.timer
[Unit]
Description=Run BetterReads pull-based deploy every 60s
Requires=bd-deploy.service

[Timer]
# Delay the first poll so reboots do not race the network coming up.
OnBootSec=2min
# Fire 60 seconds after the previous run finished, not 60s after it started.
# A slow run will not queue overlapping fires.
OnUnitInactiveSec=60s
AccuracySec=10s
# Spread fires across the minute so co-tenant timers do not align.
RandomizedDelaySec=10s
Unit=bd-deploy.service

[Install]
WantedBy=timers.target
